# crontab
# Запуск приложения каждый день в полночь (00:00)
0 0 * * * yarn start

# Запуск приложения раз в неделю в понедельник в полночь (00:00)
0 0 * * 1 yarn start
